# services/runit/transmission/run
#!/bin/zsh -l

LOGPATH="${XDG_RUNTIME_DIR}/transmission"
mkdir -p "${LOGPATH}"

exec transmission-daemon -f --log-info -e "${LOGPATH}/transmission.log"
